# script/agent/crontab
@daily docker system prune -f
@daily docker images --no-trunc --format '{{.ID}} {{.CreatedSince}}' | grep ' weeks' | awk '{ print $1 }' | xargs --no-run-if-empty docker rmi -f